# src/runtime_src/core/pcie/tools/cloud-daemon/mpd.service.in
[Unit]
Description=Xilinx Management Proxy Daemon (MPD)
Wants=network-online.target
After=network-online.target
ConditionDirectoryNotEmpty=/dev/xfpga

[Service]
Type=simple
ExecStart=${CMAKE_INSTALL_PREFIX}/${XRT_INSTALL_BIN_DIR}/mpd
Restart=on-failure

[Install]
WantedBy=multi-user.target
