# src/runtime_src/core/pcie/tools/cloud-daemon/msd.service.in
[Unit]
Description=Xilinx Management Service Daemon (MSD)
Wants=network-online.target
After=network-online.target
ConditionDirectoryNotEmpty=/dev/xfpga

[Service]
Type=simple
ExecStart=${XRT_INSTALL_DIR}/bin/msd
Restart=on-failure

[Install]
WantedBy=multi-user.target
