# src/runtime_src/tools/scripts/apu_recipes/apu-boot.service
[Unit]
Description=apu-boot
After=skd.service
 
[Service]
ExecStart=/usr/bin/apu-boot
StandardOutput=journal+console
 
[Install]
WantedBy=multi-user.target
