# src/runtime_src/tools/scripts/apu_recipes/init-apu.service
[Unit]
Description=APU init Daemon

[Service]
ExecStart=/usr/bin/init-apu
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
