# src/runtime_src/tools/scripts/apu_recipes/modprobe-vdu.service
[Unit]
Description=modprobe-vdu
 
[Service]
ExecStart=/usr/bin/modprobe-vdu
StandardOutput=journal+console
 
[Install]
WantedBy=multi-user.target
