# src/runtime_src/tools/scripts/apu_recipes/skd.service
[Unit]
Description=Soft Kernel Daemon
After=init-apu.service

[Service]
Type=forking
User=softkernel
ExecStart=/usr/bin/skd
StandardOutput=journal+console
MemoryLimit=2G

[Install]
WantedBy=multi-user.target
