# src/runtime_src/tools/scripts/apu_recipes/skd.service
[Unit]
Description=Soft Kernel Daemon
After=apu-boot.service

[Service]
Type=forking
Restart=on-failure
RestartSec=1s
User=softkernel
ExecStart=/usr/bin/skd
StandardOutput=journal+console
 
[Install]
WantedBy=multi-user.target
