# src/runtime_src/tools/scripts/apu_recipes/vdu-init.service
[Unit]
Description=vdu-init
 
[Service]
ExecStart=/usr/bin/vdu-init
StandardOutput=journal+console
 
[Install]
WantedBy=multi-user.target
